# scrapper/chiller-scraper.service
[Unit]
Description=Chiller Scraper Service
After=network.target mongod.service

[Service]
Type=simple
User=chiller
WorkingDirectory=/opt/chiller-scrapper
ExecStart=/usr/bin/node /opt/chiller-scrapper/dist/index.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production
StandardOutput=journal
StandardError=journal

[Install]
WantedBy=multi-user.target
